# misc_install_files/crontab
* * * * * /usr/bin/curl -s "http://127.0.0.1/take_sensor_observations" > /dev/null
